# timers/timertest.timer
[Unit]
Description=Runs timertest every 10s

[Timer]
OnActiveSec=2s
