# timers/timertest.service
[Unit]
Description=TimerTest

[Service]
Type=simple
ExecStart=/usr/bin/touch > /tmp/timer-test
